# crontab
0 */6 * * * cd /data/GitHub/YOUTUBE && python3 youtube_follow.py --generate-only >> /var/log/youtube-follow.log 2>&1
